# conf/supervisord/worker-download.ini
[program:worker-download]
command=sh cli php download toserver
process_name=%(program_name)s_%(process_num)02d
numprocs=10
directory=/data/www/src

autostart=true
autorestart=exited
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=QUIT
stopasgroup=false
killasgroup=false
stopwaitsecs=10

chmod = 0755
chown = webdev:webdev
user=webdev

stdout_logfile=/data/www/logs/worker-download.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/data/www/logs/worker-download.log
stderr_logfile_maxbytes=1MB
